# kakaotalk/launchd/com.jaewone.kakao-daily-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jaewone.kakao-daily-summary</string>

  <key>ProgramArguments</key>
  <array>
    <string>__ROOT__/scripts/run_daily_summary_if_needed.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>__ROOT__/logs/kakao-daily-summary.out.log</string>

  <key>StandardErrorPath</key>
  <string>__ROOT__/logs/kakao-daily-summary.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>
</dict>
</plist>
